# debian/prerm
#!/bin/bash

set -e

# Summary of how this script can be called:
# * prerm remove
# * old-prerm upgrade new-version
# * conflictor's-prerm remove in-favour package new-version
# * deconfigured's-prerm deconfigure in-favour package-being-installed version
#   [removing conflicting-package version]
# * new-prerm failed-upgrade old-version
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
    deconfigure|failed-upgrade|remove|upgrade)
        ;;
    *)
        echo "prerm called with unknown argument \`$1'" >&2
        exit 1
esac

# This deb-systemd-helper section should be generated by debhelper but
# * The debhelper tools for doing so were evolving ("dh_systemd_enable is no
#   longer used in compat >= 11, please use dh_installsystemd instead")
# * deb-systemd-helper mtime was Dec 2020 so two years unchanged
# * The template used for this deb-systemd-helper section was from 
#   https://github.com/Debian/debhelper > autoscripts > prerm-systemd with last
#   commit Jan 2018 so ~5 years unchanged

if [ -d /run/systemd/system ]; then
	deb-systemd-invoke stop sync_wireguard_to_standby.timer wireguard-logger.timer >/dev/null || true
fi

#DEBHELPER#
